# buildroot/br2-external/board/my-board/overlay/etc/profile.d/keymap.sh
#!/bin/sh
loadkeys no 2>/dev/null || true
